# packages/ositah/ositah-25.9.dev2.tar.gz/ositah-25.9.dev2/gunicorn.config/gunicorn@.service
[Install]
WantedBy=multi-user.target

[Service]
EnvironmentFile=/etc/sysconfig/gunicorn.%i
ExecStart=/bin/sh -cv '/pdisk/anaconda/envs/ositah/bin/gunicorn --chdir ${GUNICORN_APP_DIR} --pid /var/run/gunicorn.%i.pid ${GUNICORN_OPTIONS} ${GUNICORN_APP}'
PIDFile=/var/run/gunicorn.%i.pid
Restart=always
SyslogIdentifier=gunicorn

[Unit]
Description=gunicorn service %I

