# packages/django-channel-tasks/django_channel_tasks-1.1.0.dev0.tar.gz/django_channel_tasks-1.1.0.dev0/services/channel-tasks-docker.template.service
[Unit]
Description=channel-tasks

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=${CHANNEL_TASKS_HOME}/bin/restart-channel-tasks-docker.sh
Requires=docker.service
Environment="DJANGO_SECRET_KEY=${DJANGO_SECRET_KEY}"
Environment="CHANNEL_TASKS_DOCKER_HOME=${CHANNEL_TASKS_DOCKER_HOME}"
Environment="CHANNEL_TASKS_ASGI_PORT=${CHANNEL_TASKS_ASGI_PORT}"
Environment="CHANNEL_TASKS_PACKAGE=${CHANNEL_TASKS_PACKAGE}"
Environment="CHANNEL_TASKS_DB_PASSWORD=${CHANNEL_TASKS_DB_PASSWORD}"
Environment="CHANNEL_TASKS_USER=${CHANNEL_TASKS_USER}"
Environment="TASK_ADMIN_USER=${TASK_ADMIN_USER}"
Environment="TASK_ADMIN_EMAIL=${TASK_ADMIN_EMAIL}"

[Install]
WantedBy=multi-user.target
