# packages/meltingplot.duet-simplyprint-connector/meltingplot_duet_simplyprint_connector-1.3.9-py3-none-any.whl/meltingplot_duet_simplyprint_connector-1.3.9.data/data/simplyprint-connector.service
[Unit]
Description=Meltingplot Duet <-> Simplyprint.io connector
After=network.target

[Service]
Type=simple
Restart=always
StartLimitInterval=0
StartLimitBurst=1440
RestartSec=60
ExecStart=/usr/local/bin/simplyprint start
WorkingDirectory=/usr/local
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target