# packages/surf-controller/surf_controller-0.4.1.tar.gz/surf_controller-0.4.1/scheduler/crontab
# scheduler/crontab
# Run the script every night at 21:00 (9 PM)
# The script's output is redirected to logs.log inside the volume
0 21 * * * /bin/bash /app/run.sh >> /app/logs.log 2>&1

# An empty line is required at the end of the file
