# packages/django-channel-tasks/django_channel_tasks-1.1.0.dev0-py3-none-any.whl/django_channel_tasks-1.1.0.dev0.data/data/channel-tasks-services/channel-tasks.template.service
[Unit]
Description=channel-tasks

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=on-failure
ExecStart=${CHANNEL_TASKS_HOME}/bin/setup-channel-tasks.sh
Environment="DJANGO_SECRET_KEY=${DJANGO_SECRET_KEY}"
Environment="CHANNEL_TASKS_HOME=${CHANNEL_TASKS_HOME}"
Environment="CHANNEL_TASKS_SETTINGS_PATH=${CHANNEL_TASKS_SETTINGS_PATH}"
Environment="CHANNEL_TASKS_APP_NAME=${CHANNEL_TASKS_APP_NAME}"
Environment="CHANNEL_TASKS_ASGI_PORT=${CHANNEL_TASKS_ASGI_PORT}"
Environment="CHANNEL_TASKS_STATIC_ROOT=${CHANNEL_TASKS_STATIC_ROOT}"
Environment="CHANNEL_TASKS_STATIC_URI=${CHANNEL_TASKS_STATIC_URI}"
Environment="CHANNEL_TASKS_DB_PASSWORD=${CHANNEL_TASKS_DB_PASSWORD}"
Environment="CHANNEL_TASKS_USER=${CHANNEL_TASKS_USER}"
Environment="CHANNEL_TASKS_EMAIL_USER=${CHANNEL_TASKS_EMAIL_USER}"
Environment="CHANNEL_TASKS_EMAIL_PASSWORD=${CHANNEL_TASKS_EMAIL_PASSWORD}"
Environment="TASK_ADMIN_USER=${TASK_ADMIN_USER}"
Environment="TASK_ADMIN_EMAIL=${TASK_ADMIN_EMAIL}"

[Install]
WantedBy=multi-user.target
